# hal-qt.service
[Unit]
Description=Nepos Hardware abstraction layer

[Service]
ExecStart=/app/hal-qt/hal-qt
Restart=always
PrivateTmp=yes

